# image/etc/supervisor/conf.d/supervisor-apache2.conf

[program:apache2]
command=service apache2 start
autostart=true
autorestart=false
startretries=0