# image/etc/supervisor/conf.d/tomcat-apache2.conf



[program:tomcat]
command=/bin/bash -c "/usr/local/tomcat/bin/catalina.sh start"
autostart=true
autorestart=false
startretries=0
exitcodes=0
startsecs = 0
