# image/etc/supervisor/conf.d/tomcat-apache2.conf



[program:tomcat]
command=/bin/bash /usr/local/tomcat/bin/catalina.sh start
autostart=true
autorestart=false
startretries=0